# templates/default/vault.conf.erb
# This is an upstart script to keep <%= node[:vault][:server][:service_name] %> running.
#
# You can manually start and stop <%= node[:vault][:server][:service_name] %> like this:
#
# start <%= node[:vault][:server][:service_name] %>
# stop <%= node[:vault][:server][:service_name] %>
#

description "<%= node[:vault][:server][:service_name] %> service manager"

limit core unlimited unlimited

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

oom score -100
respawn

exec <%= node[:vault][:bin_dir] %>/vault server -config=<%= node[:vault][:server][:config] %> -log-level=<%= node[:vault][:server][:log_level] %>
